# infra/systemd/whisper-typer.service
[Unit]
Description=Whisper Typer RS - Speech to Text Dictation (Rust)
Documentation=https://github.com/Mizzlr/whisper-typer
After=graphical-session.target
Conflicts=whisper-typer.service code-speaker.service

[Service]
Type=simple
ExecStart=__REPO_DIR__/target/release/whisper-typer-rs
WorkingDirectory=__REPO_DIR__
Restart=on-failure
RestartSec=5

# ONNX Runtime providers (ort dlopens libonnxruntime_providers_*.so from target/release)
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus
Environment=LD_LIBRARY_PATH=__REPO_DIR__/target/release

StandardOutput=journal
StandardError=journal

# Resource limits (GPU + Kokoro ONNX)
MemoryHigh=10G
MemoryMax=14G

[Install]
WantedBy=default.target
