# infra/systemd/whisper-typer.service
[Unit]
Description=Whisper Typer RS - Speech to Text Dictation (Rust)
Documentation=https://github.com/Mizzlr/whisper-typer
After=graphical-session.target pipewire.service
PartOf=graphical-session.target
Conflicts=whisper-typer.service code-speaker.service

[Service]
Type=simple
# Free GPU memory before allocating Whisper state buffers.
ExecStartPre=-/usr/local/bin/ollama stop granite4.1:3b
ExecStart=__REPO_DIR__/target/release/whisper-typer-rs
WorkingDirectory=__REPO_DIR__
Restart=on-failure
RestartSec=5

# Environment for audio
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus
# ONNX Runtime CUDA providers; cuDNN 9 from Python venv (nvidia-cudnn-cu12 pip package)
Environment=LD_LIBRARY_PATH=__REPO_DIR__/target/release:__REPO_DIR__/.venv/lib/python3.12/site-packages/nvidia/cudnn/lib:/usr/local/lib/ollama/mlx_cuda_v13

StandardOutput=journal
StandardError=journal

# Resource limits (GPU + Kokoro ONNX)
MemoryHigh=10G
MemoryMax=14G

[Install]
WantedBy=graphical-session.target
