# systemd/whisper-typer.service
[Unit]
Description=Whisper Typer - Speech to Text Dictation
Documentation=https://github.com/Mizzlr/whisper-typer
After=graphical-session.target

[Service]
Type=simple
ExecStart=/home/mizzlr/whisper-typer/.venv/bin/python -m whisper_typer
WorkingDirectory=/home/mizzlr/whisper-typer
Restart=on-failure
RestartSec=5

# Environment for display and audio
Environment=DISPLAY=:0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus

# Logging
StandardOutput=journal
StandardError=journal

# Resource limits (adjust based on your GPU VRAM)
MemoryHigh=8G
MemoryMax=12G

[Install]
WantedBy=default.target
